# deploy/hisense-bridge.service
[Unit]
Description=Hisense MQTT Bridge for Savant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/hisense-bridge/venv/bin/python /opt/hisense-bridge/src/bridge.py /opt/hisense-bridge/config.json
WorkingDirectory=/opt/hisense-bridge
Restart=always
RestartSec=5
User=hisense-bridge

[Install]
WantedBy=multi-user.target
